# kg-api.service
[Unit]
Description=Knowledge Graph API Service
After=network.target neo4j.service redis-server.service
Wants=neo4j.service redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/knowledge-graph/api
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/kg-api.log
StandardError=append:/var/log/kg-api-error.log

[Install]
WantedBy=multi-user.target

